# systemd/skynet-v18-fade-db-shadow.service.d/99-retired.conf
[Unit]
# The obsolete low-cost Phase2D service may run only after
# an explicit manual creation of this opt-in file.
ConditionPathExists=/root/skynet/ALLOW_RETIRED_V18_FADE
