# systemd/skynet-v19-preregistered-forward.timer.d/99-retired.conf
[Unit]
ConditionPathExists=/root/skynet/ALLOW_REJECTED_V19_FORWARD_RECHECK
