# systemd/skynet-v18-exact-multi-shadow.service.d/10-wallpos-selftest.conf
[Service]
ExecStartPre=/root/skynet/.venv/bin/python /root/skynet/v18_exact_wallpos_selftest.py
