# systemd/skynet-v19-preregistered-forward.timer
[Unit]
Description=Evaluate SKYNET V19 preregistered forward gates every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=skynet-v19-preregistered-forward.service

[Install]
WantedBy=timers.target
