# systemd/skynet-v18-exact-forward-gate.timer
[Unit]
Description=Run SKYNET V18 exact forward gate every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=skynet-v18-exact-forward-gate.service

[Install]
WantedBy=timers.target
