# systemd/skynet-v18-archive-walkforward.service
[Unit]
Description=SKYNET V18 archive census and fixed WALL_POS walkforward
After=skynet-recorder.service skynet-v18-exact-multi-shadow.service

[Service]
Type=oneshot
WorkingDirectory=/root/skynet
ExecStart=/root/skynet/.venv/bin/python /root/skynet/v18_archive_walkforward_lab.py
User=root
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=infinity
Environment=PYTHONUNBUFFERED=1
